# paperless-ai.container
[Unit]
Description=Paperless AI — AI companion for Paperless-ngx
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/clusterzx/paperless-ai:latest
AutoUpdate=registry

# The image runs as root (uid 0) inside the container.
# In rootless Podman the default user-namespace mapping makes root-in-container
# equivalent to the service user on the host, so no UserNS=keep-id is needed.

# Environment variables
EnvironmentFile=%h/.config/containers/systemd/paperless-ai.env
EnvironmentFile=%h/.config/containers/systemd/paperless-ai.override.env

# Persistent data volume — /app/data holds the SQLite database and processed-document state
Volume=%h/data:/app/data:Z

# Port mapping — bind to localhost only, reverse proxy handles TLS
PublishPort=127.0.0.1:3000:3000

# Health check
HealthCmd=curl -f http://localhost:3000/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
