# paperless-ai-backup.timer
[Unit]
Description=Daily Paperless AI backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
